# rootfs/etc/services.d/squid/run
#!/usr/bin/execlineb -P
/usr/sbin/squid -N -d1
